# co.invariante.zeit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Label: unique identifier for this job -->
    <key>Label</key>
    <string>co.invariante.zeit</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/miguelvilagonzalez/repos/zeit/run_tracker.sh</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/miguelvilagonzalez/repos/zeit</string>

    <!-- Run every 60 seconds -->
    <key>StartInterval</key>
    <integer>60</integer>

    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/miguelvilagonzalez/repos/zeit/logs/launchd.out.log</string>

    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/Users/miguelvilagonzalez/repos/zeit/logs/launchd.err.log</string>

    <!-- Run even if no user is logged in (optional, remove if you only want it when logged in) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep the job alive - restart if it crashes (optional) -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
