# co.invariante.zeit.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Label: unique identifier for this job -->
    <key>Label</key>
    <string>co.invariante.zeit.menubar</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/miguelvilagonzalez/repos/zeit/.venv/bin/python</string>
        <string>/Users/miguelvilagonzalez/repos/zeit/run_menubar_app.py</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/miguelvilagonzalez/repos/zeit</string>

    <!-- Start on load (when user logs in) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Only restart on crash (non-zero exit), not on clean Quit -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/miguelvilagonzalez/repos/zeit/logs/menubar.out.log</string>

    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/Users/miguelvilagonzalez/repos/zeit/logs/menubar.err.log</string>

    <!-- Environment variables (optional, but good to set) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Process type - this is important for GUI apps -->
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
